# AI AGENT/watchdog/openclaw-watchdog.timer
[Unit]
Description=OpenClaw Watchdog Timer (every 30s)

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
