# AI AGENT/watchdog/openclaw-watchdog.service
[Unit]
Description=OpenClaw Watchdog Health Check
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.openclaw/watchdog/watchdog.sh
WorkingDirectory=%h/.openclaw

[Install]
WantedBy=default.target
